# hapro/rootfs/etc/s6-overlay/s6-rc.d/setup-addon/run
#!/usr/bin/with-contenv bashio
bashio::log.info "SETUP-ADDON: Addon version is $(bashio::addon.version)"
versionAvailable=$(bashio::addon.update_available)
if [ "$versionAvailable" == "true" ]; then
    bashio::log.warning "SETUP-ADDON: New addon version is available $(bashio::addon.version_latest)"
fi
bashio::log.info "SETUP-ADDON: Starting setup"

# Load request helper functions
source /usr/src/request.sh

# Get API URL from config
HaproApi=$(bashio::config 'apiUrl')
HaproApi="${HaproApi%/}" # Remove trailing slash if exists
bashio::log.info "SETUP-ADDON: Using API URL: $HaproApi"

bashio::log.notice "SETUP-ADDON: Setting HA configuration"
mkdir -p /homeassistant/custom_components
if [ ! -d "/homeassistant/custom_components/auth_header" ]; then
    cp -R /usr/src/hass-auth-header/custom_components/* /homeassistant/custom_components/
fi
if [ ! -f "/homeassistant/hapro-http.yaml" ]; then
    cp /usr/src/config/http.yaml /homeassistant/hapro-http.yaml
fi
if [ ! -f "/homeassistant/hapro-auth.yaml" ]; then
    cp /usr/src/config/auth.yaml /homeassistant/hapro-auth.yaml
fi

changed=0
if ! grep -q "http:" /homeassistant/configuration.yaml; then
    echo "http: !include hapro-http.yaml #Added by the HaPro addon" >> /homeassistant/configuration.yaml
    changed=1
else
    if ! grep -q "^http: !include hapro-http.yaml #Added by the HaPro addon" /homeassistant/configuration.yaml; then
        http_section=$(sed -n '/^http:/,/^[^[:space:]-]/ { /^http:/b; p }' "/homeassistant/configuration.yaml")
        use_x_forwarded_for=$(grep -E "^use_x_forwarded_for:" "/homeassistant/hapro-http.yaml")
        if ! echo "$http_section" | grep -q "use_x_forwarded_for:"; then
            sed -i '/^http:/a \ \ use_x_forwarded_for: true #added by hapro' "/homeassistant/configuration.yaml"
            changed=1
        elif echo "$http_section" | grep -q "use_x_forwarded_for: false"; then
            sed -i 's/use_x_forwarded_for: false/use_x_forwarded_for: true #modified by hapro/' "/homeassistant/configuration.yaml"
            changed=1
        fi
       if echo "$http_section" | grep -q "trusted_proxies:"; then
            for ip in "::0" "::1" "0.0.0.0" "127.0.0.1" "172.23.0.0/24" "172.17.0.1" "172.17.0.0/24"; do
                if ! echo "$http_section" | grep -q "  - $ip"; then
                    sed -i "/trusted_proxies:/a \ \ \ \ - $ip #added by hapro" "/homeassistant/configuration.yaml"
                    changed=1
                fi
            done
        else
            sed -i '/^http:/a \ \ trusted_proxies:\n\ \ \ \ - ::0 #added by hapro\n\ \ \ \ - ::1 #added by hapro\n\ \ \ \ - 0.0.0.0 #added by hapro\n\ \ \ \ - 127.0.0.1 #added by hapro\n\ \ \ \ - 172.23.0.0/24 #added by hapro\n\ \ \ \ - 172.17.0.1 #added by hapro\n\ \ \ \ - 172.17.0.0/24 #added by hapro' "/homeassistant/configuration.yaml"
            changed=1
        fi
    fi
fi

if grep -q "^[[:space:]]*#\+.*auth_header: !include hapro-auth.yaml" /homeassistant/configuration.yaml; then # present but disabled
    sed -i 's/^[[:space:]]*#\+[[:space:]]*auth_header:.*$/auth_header: !include hapro-auth.yaml #Added by the HaPro addon/' /homeassistant/configuration.yaml
    changed=1
elif ! grep -q "^[[:space:]]*auth_header: !include hapro-auth.yaml" /homeassistant/configuration.yaml; then # Not present at all, add it
    echo "auth_header: !include hapro-auth.yaml #Added by the HaPro addon" >> /homeassistant/configuration.yaml
    changed=1
fi
if [ "$changed" -eq 1 ]; then
    bashio::log.error "SETUP-ADDON: Made changes to configuration.yaml, please restart Home Assistant"
    doHaInternalApiRequest "services/notify/persistent_notification" "POST" "{\"title\":\"HaPro requires restart\",\"message\":\"HaPro changed your configuration.yaml, please restart Home Assistant\",\"data\":{\"notification_id\":\"hapro_restart\"}}"
    /etc/s6-overlay/s6-rc.d/setup-addon/finish 1
    exit 1
fi

bashio::log.notice "SETUP-ADDON: Getting client.toml"
uuid=$(bashio::config 'uuid')
bashio::log.notice "SETUP-ADDON: Current UUID: $uuid"
if [ -z "$uuid" ] || [ "$uuid" == "null" ]; then
    uuid=$(cat /proc/sys/kernel/random/uuid | tr -d '-')
    bashio::log.info "SETUP-ADDON: No UUID set in addon config, generated new UUID: $uuid"
    bashio::addon.option "uuid" "$uuid"
fi

res=$(curl -k -s -o /usr/bin/client.toml -w "%{http_code}" "$HaproApi/api/Device/$uuid/rathole" )
if [ "$res" -eq 404 ]; then
    bashio::log.error "SETUP-ADDON: Failed to get client.toml, registering device first"

    bashio::log.notice "SETUP-ADDON: Getting owner account"
    # doSupervisorRequest "auth/list"
    ownerAcc=$(jq -r '.data.users[] | select(.is_owner==true) | .name' /homeassistant/.storage/auth)
    bashio::log.notice "SETUP-ADDON: Owner account: $ownerAcc"

    res2=$(curl -k -s -H "Content-Type: application/json" -d "{\"uuid\":\"$uuid\",\"haAcc\":\"$ownerAcc\"}" "$HaproApi/api/PendingDevice/"| jq .statusCode)
    bashio::log.notice "SETUP-ADDON: HTTP Post device Status Code: $res2"
    if [ "$res2" -eq 200 ]; then
        bashio::log.info "SETUP-ADDON: Device registered, claim this device on the portal with token: $uuid"
    fi
    /etc/s6-overlay/s6-rc.d/setup-addon/finish 1
    exit 1
elif [ "$res" -ne 200 ]; then
    bashio::log.error "SETUP-ADDON: HTTP Get Client.toml Status Code: $res"
    bashio::log.error "SETUP-ADDON: Something went wrong, please try again later"
    /etc/s6-overlay/s6-rc.d/setup-addon/finish 1
    exit 1
fi

if [ ! -f "/homeassistant/blueprints/automation/hapro_notification.yaml" ]; then
    cp /usr/src/config/hapro_notification.yaml /homeassistant/blueprints/automation/hapro_notification.yaml
    doHaInternalApiRequest "services/automation/reload" "POST"
elif ! cmp -s /homeassistant/blueprints/automation/hapro_notification.yaml /usr/src/config/hapro_notification.yaml; then
    cp /usr/src/config/hapro_notification.yaml /homeassistant/blueprints/automation/hapro_notification.yaml
    doHaInternalApiRequest "services/automation/reload" "POST"
fi


bashio::log.info "SETUP-ADDON: Setup done"