# hapro/rootfs/etc/s6-overlay/s6-rc.d/rathole/run
#!/usr/bin/with-contenv bashio
source /usr/src/request.sh

caddy_port=$(bashio::config 'caddy_port')
sed -i "s/{{.PORT}}/$caddy_port/g" /usr/bin/client.toml
sed -i "s/{{.HOST}}/127.0.0.1/g" /usr/bin/client.toml

bashio::log.info "RATHOLE: Starting rathole"
if bashio::config.true 'debug'; then
    export RUST_LOG=debug
    bashio::log.notice "RATHOLE: Debug mode enabled"
fi
exec /usr/bin/rathole -c /usr/bin/client.toml  2>&1 | while read -r line; do
    case "$line" in
        *" DEBUG "*) bashio::log.notice "RATHOLE: ${line#* DEBUG }" ;;
        *" INFO "*) bashio::log.info "RATHOLE: ${line#* INFO }" ;;
        *" WARN "*) bashio::log.warning "RATHOLE: ${line#* WARN }" ;;
        *" ERROR "*) bashio::log.error "RATHOLE: ${line#* ERROR }" ;;
        *) bashio::log.notice "RATHOLE: $line" ;;
    esac
done